# maint-scripts/nunet-dms/etc/systemd/system/nomad-client-restart.path
[Path]
PathChanged=/etc/nunet/clientV2.json

[Install]
WantedBy=multi-user.target
